# .ebextensions/crontab
1 18 * * * sudo sh /var/app/current/.ebextensions/run_all_docker.sh
